# uppoint-cloud/ops/systemd/uppoint-cloud.service
[Unit]
Description=Uppoint Cloud Next.js Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/uppoint-cloud
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=-/etc/uppoint-cloud.env
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5
TimeoutStartSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/uppoint-cloud/.next

[Install]
WantedBy=multi-user.target
