# uppoint-cloud/ops/systemd/uppoint-cloud.service
[Unit]
Description=Uppoint Cloud Next.js Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/uppoint-cloud
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/opt/uppoint-cloud/.env
ExecStartPre=+/usr/bin/install -d -o www-data -g www-data -m 0755 /opt/uppoint-cloud/.next
ExecStartPre=+/usr/bin/install -d -o www-data -g www-data -m 0755 /opt/uppoint-cloud/.next/cache
ExecStartPre=+/usr/bin/install -d -o www-data -g www-data -m 0755 /opt/uppoint-cloud/.next/cache/images
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5
TimeoutStartSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/uppoint-cloud/.next /opt/uppoint-cloud/.next/cache /opt/uppoint-cloud/.next/cache/images

[Install]
WantedBy=multi-user.target
