# uppoint-cloud/ops/systemd/uppoint-cloud.service
[Unit]
Description=Uppoint Cloud Next.js Service
After=network-online.target uppoint-tune.service
Wants=network-online.target uppoint-tune.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/uppoint-cloud
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/opt/uppoint-cloud/.env
# tune-system.sh tarafından oluşturulur; yoksa yoksay (-).
EnvironmentFile=-/etc/uppoint-cloud-tuned.env
ExecStartPre=+/usr/bin/install -d -o www-data -g www-data -m 0755 /opt/uppoint-cloud/.next
ExecStartPre=+/usr/bin/install -d -o www-data -g www-data -m 0755 /opt/uppoint-cloud/.next/cache
ExecStartPre=+/usr/bin/install -d -o www-data -g www-data -m 0755 /opt/uppoint-cloud/.next/cache/images
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15
# RAM'in %20'si — donanım yükseltilirse otomatik olarak ölçeklenir.
MemoryMax=20%
# Dosya tanımlayıcı üst sınırı — her bağlantı + statik dosya için gerekli.
LimitNOFILE=65536
# CPU kullanımını 4 çekirdekle sınırla (8 çekirdekli sunucu); ani CPU spike'larını engeller.
CPUQuota=400%
# Süreç/iş parçacığı limiti — fork bomb koruması.
TasksMax=512
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/uppoint-cloud/.next /opt/uppoint-cloud/.next/cache /opt/uppoint-cloud/.next/cache/images

[Install]
WantedBy=multi-user.target
